# debian/apt-daily.service
[Unit]
Description=Daily apt activities
Documentation=man:apt(8)

[Service]
Type=oneshot
ExecStart=/usr/lib/apt/apt.systemd.daily

